# envs/ss4.9/php7.4/docker/crontab
*/1 * * * * /var/www/html/vendor/bin/sake dev/tasks/ProcessJobQueueTask
